# crontab/crons
*/5 * * * * /home/walter/.local/bin/batterynotify
